# launchd/com.pet-care.local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pet-care.local</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="/Users/lyl168/.nvm/versions/node/v24.14.1/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin"; cd "/Users/lyl168/Desktop/pet_care"; npm run start -- -p 3000</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/lyl168/Desktop/pet_care</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/lyl168/Desktop/pet_care/logs/local-server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/lyl168/Desktop/pet_care/logs/local-server.err.log</string>
</dict>
</plist>
